# scripts/hermes-hub-linux.service
[Unit]
Description=Hermes Gateway API for Hermes Hub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=HERMES_API_HOST=0.0.0.0
Environment=HERMES_API_PORT=8642
Environment=API_SERVER_KEY=hermes-hub
Environment=HERMES_MAX_ITERATIONS=0
Environment=HERMES_AUXILIARY_LOCAL_ONLY=true
Environment=HERMES_INFERENCE_PROVIDER=lm_studio
Environment=LM_STUDIO_BASE_URL=http://127.0.0.1:1234
Environment=HERMES_VIDEO_LIBRARY_PATH=%h/.hermes/media/video
Environment=HERMES_MEDIA_ROOTS=%h:%h/.hermes/cache:%h/.hermes/media:%h/.hermes/media/video
Environment=HERMES_HUB_STATE_PATH=%h/.hermes/hub_state.json
Environment=HERMES_HUB_MEMORY_PATH=%h/.hermes/hub_memory.json
# Final Ubuntu headless/vLLM profile:
# Environment=HERMES_INFERENCE_PROVIDER=vllm
# Environment=VLLM_BASE_URL=http://127.0.0.1:8000
# Environment=HERMES_INFERENCE_MODEL=your-vllm-model-id
ExecStart=%h/hermes-hub-linux.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
