# scripts/hermes-hub-linux.service
[Unit]
Description=Hermes Gateway API for Hermes Hub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PATH=/home/matteo/.local/bin:/home/matteo/.hermes/bin:/home/matteo/.hermes/node/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HERMES_API_HOST=0.0.0.0
Environment=HERMES_API_PORT=8642
Environment=API_SERVER_HOST=0.0.0.0
Environment=API_SERVER_PORT=8642
Environment=API_SERVER_KEY=hermes-hub
Environment=HERMES_API_KEY=hermes-hub
Environment=HERMESAPIKEY=hermes-hub
Environment=HERMES_HUB_API_KEY=hermes-hub
Environment=HERMES_GATEWAY_API_KEY=hermes-hub
Environment=HERMES_MAX_ITERATIONS=120
Environment=HERMES_AUXILIARY_LOCAL_ONLY=true
Environment=HERMES_INFERENCE_PROVIDER=custom
Environment=HERMES_INFERENCE_BASE_URL=http://127.0.0.1:8000/v1
Environment=HERMES_INFERENCE_MODEL=HauhauCS/Qwen3.6-35B-A3B-Uncensored-HauhauCS-Aggressive:IQ4_XS
Environment=HERMES_VIDEO_LIBRARY_PATH=%h/.hermes/media/video
Environment=HERMES_MEDIA_ROOTS=%h:%h/.hermes/cache:%h/.hermes/media:%h/.hermes/media/video
Environment=HERMES_HUB_STATE_PATH=%h/.hermes/hub_state.json
Environment=HERMES_HUB_MEMORY_PATH=%h/.hermes/hub_memory.json
ExecStartPre=/home/matteo/.local/bin/hermes-wait-tailscale
ExecStartPre=/home/matteo/.local/bin/hermes-wait-llama
ExecStart=%h/hermes-hub-linux.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=1000
TimeoutStopSec=240

[Install]
WantedBy=default.target
