# scripts/hermes-hub-linux-update.service
[Unit]
Description=Update Hermes Hub Linux gateway helper
Documentation=https://github.com/JackoPeru/HermesHub

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hermes-hub-linux-update --restart
# API paginata (fino a 5 x 60 s) + download (900 s) + restart gateway (1300 s) + probe + margine.
TimeoutStartSec=3000
