# scripts/hermes-hub-linux.service
[Unit]
Description=Hermes Hub API Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=HERMES_API_HOST=0.0.0.0
Environment=HERMES_API_PORT=8642
Environment=LM_STUDIO_BASE_URL=http://127.0.0.1:1234
ExecStart=%h/hermes-hub-linux.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
